# systemd/user/battery-check.service
[Unit]
Description=Check battery level and notify when charging

[Service]
ExecStart=/home/yas/.config/i3/scripts/battery-check.sh

